# deploy/ilay-frontend.service
[Unit]
Description=ILAY Next.js Frontend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ilay/frontend
Environment=NODE_ENV=production
Environment=NEXT_PUBLIC_API_URL=http://134.209.228.120
ExecStart=/usr/bin/npx next start --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
